# docker/root/etc/s6-overlay/s6-rc.d/svc-workers/run
#!/command/with-contenv sh
# shellcheck shell=bash

cd /app/apps/workers;
if [ "${USE_JEMALLOC:-false}" = "true" ]; then
    export LD_PRELOAD=/usr/local/lib/libjemalloc.so
fi
exec node dist/index.mjs
